# services/stt/quadlet/bloom-stt.container
[Unit]
Description=Bloom STT — Speech-to-text (whisper.cpp server, OpenAI-compatible)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/ggml-org/whisper.cpp:main
ContainerName=bloom-stt

# Bridge network for isolation
Network=bloom.network

# Expose API (8080 inside -> 8081 outside, localhost + NetBird mesh, firewall-protected)
PublishPort=0.0.0.0:8081:8080

# Model storage persists across restarts
Volume=bloom-stt-models:/models

# Media files for transcription (read-only)
Volume=/var/lib/bloom/media:/media:ro

# Override bash entrypoint with actual binary
PodmanArgs=--entrypoint /app/build/bin/whisper-server

# Run whisper server with base model
Exec=--host 0.0.0.0 --port 8080 --model /models/ggml-base.en.bin --threads 4

PodmanArgs=--memory=1g
PodmanArgs=--security-opt label=disable
HealthCmd=curl -sf http://localhost:8080/health || exit 1
HealthInterval=30s
HealthRetries=3
HealthTimeout=10s
HealthStartPeriod=2m
NoNewPrivileges=true
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=default.target
