# core/os/system_files/usr/lib/systemd/system/bloom-matrix.service
[Unit]
Description=Bloom Matrix Homeserver (Continuwuity)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Generate registration token on first boot if missing
ExecStartPre=/bin/bash -c 'TOKEN_FILE=/var/lib/continuwuity/registration_token; [ -f "$TOKEN_FILE" ] || (openssl rand -base64 32 > "$TOKEN_FILE" && chmod 640 "$TOKEN_FILE")'
ExecStart=/usr/local/bin/continuwuity
Environment=CONTINUWUITY_CONFIG=/etc/bloom/matrix.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
DynamicUser=yes
StateDirectory=continuwuity
RuntimeDirectory=continuwuity

[Install]
WantedBy=multi-user.target
