# core/os/system_files/usr/lib/systemd/system/bloom-update-check.service
[Unit]
Description=Check for Bloom OS updates
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
ExecStart=/usr/local/bin/bloom-update-check.sh
