# services/syncthing/quadlet/bloom-syncthing.container
[Unit]
Description=Bloom Syncthing — Home directory sync service
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/syncthing/syncthing@sha256:1feffa2d4826b48f25faefed093d07c5f00304d7e7ac86fd7cda334d22651643
ContainerName=bloom-syncthing
Network=host

# Persistent Syncthing config + index state
Volume=bloom-syncthing-data:/var/syncthing:Z

# Sync target: user home directory
Volume=%h:/var/syncthing/home:Z

Environment=PUID=%U
Environment=PGID=%G
Environment=STGUIADDRESS=0.0.0.0:8384
PodmanArgs=--memory=256m

# Health endpoint does not require API key and is stable for container checks
HealthCmd=curl -sf http://127.0.0.1:8384/rest/noauth/health > /dev/null || exit 1
HealthInterval=30s
HealthRetries=3
HealthTimeout=10s

NoNewPrivileges=true
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=default.target
