# core/os/system_files/usr/lib/systemd/system/netbird.service.d/10-bloom.conf
[Unit]
# Ensure /var/log/netbird exists before starting (created by tmpfiles.d/bloom.conf)
After=systemd-tmpfiles-setup.service
