# os/sysconfig/bloom-sway.service
[Unit]
Description=Bloom Display (Sway Wayland compositor — headless VNC mode)
# Only start if no DRM GPU is available (headless server mode).
# When a GPU is present, Sway starts from .bash_profile on tty1 instead,
# which gives it a proper logind session for DRM/VT access.
ConditionPathIsDirectory=!/dev/dri
After=network.target

[Service]
Type=simple
User=pi
RuntimeDirectory=bloom
Environment=XDG_RUNTIME_DIR=/run/user/%U
ExecStart=/usr/local/share/bloom/os/scripts/start-sway.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=bloom-display.target
