# os/sysconfig/bloom-display.service
[Unit]
Description=Bloom Display (Xpra + Xvfb)
After=network.target

[Service]
User=bloom
Environment=DISPLAY=:99
ExecStart=/usr/bin/xpra start :99 \
    --opengl=no \
    --bind-tcp=0.0.0.0:14500 \
    --html=on \
    --no-daemon
ExecStop=/usr/bin/xpra stop :99
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
