# os/sysconfig/bloom-llm-local.service
[Unit]
Description=Bloom Local LLM (llama.cpp + Crow-4B Opus Distill)
After=network.target
Before=getty@tty1.service

[Service]
Type=simple
User=llm
Group=llm
ExecStart=/usr/local/bin/llama-server \
    --model /usr/local/share/bloom/models/crow-4b.gguf \
    --host 127.0.0.1 \
    --port 8080 \
    --ctx-size 8192 \
    --threads 4 \
    --no-mmap \
    --chat-template chatml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
