# os/sysconfig/bloom-display.service
[Unit]
Description=Bloom Display (Xvfb headless framebuffer)
After=network.target

[Service]
User=bloom
Environment=DISPLAY=:99
ExecStart=/usr/bin/Xvfb :99 -screen 0 1280x1024x24
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
