# services/code-server/quadlet/bloom-code-server.container
[Unit]
Description=Bloom code-server — web-based code editor
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/bloom-code-server:latest
ContainerName=bloom-code-server
# Host networking so code-server is directly accessible via NetBird mesh (no port forwarding needed)
Network=host
Volume=%h:/home/coder/project
Environment=DISABLE_TELEMETRY=true
PodmanArgs=--memory=512m
PodmanArgs=--security-opt label=disable
NoNewPrivileges=true
LogDriver=journald

HealthCmd=curl -f http://localhost:8443/healthz || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=10s
HealthTimeout=5s

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=default.target
