# services/matrix/quadlet/bloom-matrix.container
[Unit]
Description=Bloom Matrix Homeserver (Continuwuity)
After=network-online.target
Wants=network-online.target

[Container]
Image=forgejo.ellis.link/continuwuation/continuwuity:latest
ContainerName=bloom-matrix
Network=bloom.network
PublishPort=6167:6167
Volume=bloom-matrix-data:/var/lib/continuwuity

Environment=CONTINUWUITY_SERVER_NAME=bloom
Environment=CONTINUWUITY_DATABASE_PATH=/var/lib/continuwuity
Environment=CONTINUWUITY_PORT=6167
Environment=CONTINUWUITY_ALLOW_FEDERATION=false
Environment=CONTINUWUITY_ALLOW_REGISTRATION=true
Environment=CONTINUWUITY_ADDRESS=0.0.0.0
Environment=CONTINUWUITY_MAX_REQUEST_SIZE=20000000
Environment=CONTINUWUITY_ALLOW_CHECK_FOR_UPDATES=false

EnvironmentFile=%h/.config/bloom/matrix.env

PodmanArgs=--memory=256m
HealthCmd=wget -qO- http://localhost:6167/_matrix/client/versions || exit 1
HealthInterval=30s
HealthRetries=3
HealthTimeout=10s
HealthStartPeriod=30s
NoNewPrivileges=true
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=default.target
