# services/_template/quadlet/bloom-TEMPLATE.container
[Unit]
Description=Bloom TEMPLATE Service
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/bloom-TEMPLATE:latest
ContainerName=bloom-TEMPLATE
User=0

# Host networking for direct mesh access
Network=host

# Service state persists across restarts
# TODO: Adjust volume name and mount path for your service
Volume=bloom-TEMPLATE-data:/data/TEMPLATE

Environment=NODE_ENV=production

# Service-specific config
EnvironmentFile=%h/.config/bloom/TEMPLATE.env

PodmanArgs=--memory=256m
PodmanArgs=--security-opt label=disable
# TODO: Assign a unique port (check catalog.yaml for used ports)
HealthCmd=wget -qO- http://localhost:18800/health || exit 1
HealthInterval=30s
HealthRetries=3
HealthTimeout=10s
HealthStartPeriod=60s
NoNewPrivileges=true
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=default.target
