# services/llm/quadlet/bloom-llm.container
[Unit]
Description=Bloom LLM — Local language model (llama.cpp server, OpenAI-compatible)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/ggml-org/llama.cpp:server
ContainerName=bloom-llm

# Bridge network for isolation
Network=bloom.network

# Expose OpenAI-compatible API (localhost + NetBird mesh, firewall-protected)
PublishPort=0.0.0.0:8080:8080

# Model storage persists across restarts
Volume=bloom-llm-models:/models

# Run server with model (downloaded by first-boot skill)
Exec=--host 0.0.0.0 --port 8080 --model /models/default.gguf --ctx-size 2048 --threads 4

PodmanArgs=--memory=3g
PodmanArgs=--security-opt label=disable
HealthCmd=curl -sf http://localhost:8080/health || exit 1
HealthInterval=30s
HealthRetries=3
HealthTimeout=10s
HealthStartPeriod=2m
NoNewPrivileges=true
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=default.target
