# services/lemonade/quadlet/bloom-lemonade.container
[Unit]
Description=Bloom Lemonade — Local LLM + speech-to-text (OpenAI-compatible API)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/lemonade-sdk/lemonade-server:latest
ContainerName=bloom-lemonade

# Bridge network for isolation
Network=bloom.network

# Expose OpenAI-compatible API on localhost
PublishPort=127.0.0.1:8000:8000

# Model cache persists across restarts
Volume=bloom-lemonade-models:/root/.cache/huggingface

# Media files for transcription (read-only)
Volume=/var/lib/bloom/media:/media:ro

Environment=LEMONADE_LLAMACPP_BACKEND=cpu
PodmanArgs=--memory=4g
PodmanArgs=--security-opt label=disable
HealthCmd=curl -sf http://localhost:8000/health || exit 1
HealthInterval=30s
HealthRetries=3
HealthTimeout=10s
NoNewPrivileges=true
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=600
TimeoutStopSec=30

[Install]
WantedBy=default.target
